# deps/manager-api.service
[Unit]
Description=APACHE APISIX Manager-api Service
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/manager-api -p /usr/local/apisix/dashboard

[Install]
WantedBy=multi-user.target
